# .root/etc/udev/rules.d/90-devices.rules
KERNEL=="card0", SUBSYSTEM=="drm", ENV{DISPLAY}=":0", ENV{XAUTHORITY}="/home/zoresvit/.Xauthority", RUN+="/usr/local/bin/displays.sh"
